# files/etc/skel/.config/autostart/first-login.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=first-login.sh
Exec=/usr/local/bin/first-login.sh
Terminal=false
